# configs/supervisor/openclaw-gateway.conf
[program:openclaw-gateway]
environment=PROC_NAME="%(program_name)s",OLLAMA_API_KEY="ollama-local"
command=/opt/supervisor-scripts/openclaw-gateway.sh
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10
priority=200
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_events_enabled=true
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
